# dot_config/systemd/user/offlinemsmtp.service
[Unit]
Description=Offline msmtp

[Service]
ExecStart=/usr/bin/offlinemsmtp --daemon
Restart=always
RestartSec=0

[Install]
WantedBy=default.target
