# dot_config/systemd/user/wallpaper.service
[Unit]
Description=Set wallpaper

[Service]
ExecStart=/home/sumner/bin/set_wallpaper.sh

[Install]
WantedBy=default.target
